# bin/raspi-temperature-server.service
[Unit]
Description=Start the raspi-temperature-server meteor app
After=network.target

[Service]
ExecStart=/usr/bin/meteor
WorkingDirectory=/home/pi/raspi-temperature-server
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
